# pdf-ms.service
[Unit]
Description=PDF Split & Merge
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/path/to/pdf_ms
ExecStart=/path/to/venv/bin/python api/server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
